# root/etc/s6-overlay/s6-rc.d/init-config-jcdssync/run
#!/usr/bin/with-contenv bash

mkdir -p /packages

# permissions
chown -R abc:abc \
	/config
chown abc:abc \
	/packages
